# package/base-files/files/etc/init.d/network-watcher
#!/bin/sh /etc/rc.common
START=95  # 启动顺序（晚于network/ModemManager）
STOP=01

USE_PROCD=1

start_service() {
    procd_open_instance
    procd_set_param command /bin/sh -c "
        while true; do
            # 检查5G（wwan0）状态：无IP则重启ModemManager
            if ! ip addr show wwan0 | grep -q 'inet '; then
                logger -t NetworkWatcher '5G (wwan0) disconnected, restarting ModemManager...'
                /etc/init.d/ModemManager restart
                sleep 30
            fi

            # 检查WIFI（wlan0）状态：无接口则重启hostapd
            if ! ip link show wlan0 | grep -q 'UP'; then
                logger -t NetworkWatcher 'WIFI (wlan0) disconnected, restarting hostapd...'
                /etc/init.d/hostapd restart
                sleep 10
            fi

            sleep 60  # 每60秒检查一次
        done
    "
    procd_set_param respawn  # 进程退出自动重启
    procd_close_instance
}

stop_service() {
    procd_stop_instance
}
ENABLE=1
