# custom-files/etc/init.d/5g-auto
#!/bin/sh /etc/rc.common

START=99
STOP=15

start() {
    echo "启动5G自动连接服务..."
    /usr/bin/5g-auto-start.sh &
}

stop() {
    echo "停止5G自动连接服务..."
    pkill -f "5g-auto-start.sh"
}
